# breaktimer-core.service
[Unit]
Description=breaktimer core — headless work-session timer
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/python3 /home/user/Projects/breaktimer/main.py
WorkingDirectory=/home/user/Projects/breaktimer
Restart=always
RestartSec=5
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/etc/hosts
NoNewPrivileges=true

[Install]
WantedBy=default.target
